# deploy/systemd/transcription-calls-911.service
[Unit]
Description=Enqueue 911 weekly pipeline
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/transcription-calls
# docker compose сам прочитает .env (и .image.env через --env-file).
ExecStart=/usr/bin/docker compose --env-file .env --env-file .image.env run --rm web python scripts/enqueue.py n911_run
