# deploy/systemd/transcription-calls-classify.timer
[Unit]
Description=Periodic classification enqueuer

[Timer]
OnCalendar=*:0/10
Persistent=true
Unit=transcription-calls-classify.service

[Install]
WantedBy=timers.target

