# deploy/systemd/transcription-calls-summarize911.service
[Unit]
Description=Enqueue pending 911 summarization
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/transcription-calls
EnvironmentFile=/opt/transcription-calls/prod.env
ExecStart=/usr/bin/docker compose run --rm web python scripts/enqueue.py summarize_911 --limit 100

