# deploy/systemd/transcription-calls-kc.service
[Unit]
Description=Enqueue KC daily pipeline
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/transcription-calls
ExecStart=/usr/bin/docker compose --env-file .env --env-file .image.env run --rm web python scripts/enqueue.py kc_day
